# etc/supervisor/conf.d/rsyslog.conf
[program:rsyslog]
command=/bin/bash -c "rsyslogd"
